# dolphin.service
[Unit]
Description=Integration Manager

[Service]
ExecStart=~/dolphin/run.sh

[Install]
WantedBy=multi-user.target