# dolphin.service
[Unit]
Description=Integration Manager

[Service]
ExecStart=./usr/local/src/dolphin/run.sh

[Install]
WantedBy=multi-user.target